# src/cursor/root/etc/s6-overlay/s6-rc.d/svc-cursor/run
#!/usr/bin/with-contenv bash

# Start Cursor IDE
exec s6-setuidgid abc /Cursor.AppImage --no-sandbox --appimage-extract-and-run
