# src/cursor/root/etc/s6-overlay/s6-rc.d/svc-cursor/run
#!/usr/bin/with-contenv bash

# Wait for KasmVNC X display to be ready
until xdpyinfo -display :1 >/dev/null 2>&1; do sleep 1; done

# Start Cursor IDE
exec s6-setuidgid abc \
    env DISPLAY=:1 \
    /Cursor.AppImage --no-sandbox --appimage-extract-and-run
